# meta-kos-core/recipes-core/custom-boot/files/knowledge-splash.service
[Unit]
Description=Knowledge OS Splash
DefaultDependencies=no
After=local-fs.target
Before=getty@tty1.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'clear > /dev/tty1; echo -e "\033[1;32mWelcome to Knowledge OS\033[0m" > /dev/tty1'
StandardOutput=tty
TTYPath=/dev/tty1

[Install]
WantedBy=sysinit.target